# docker/openstack/neutron/neutron-l3-agent/root/etc/service/neutron-ha-tool/run
#!/usr/bin/env sh

sleep 1
exec bash /usr/bin/neutron-ha-tool.sh

