# docker/openstack/neutron/neutron-server/root/etc/service/neutron-server/run
#!/usr/bin/env sh

exec /usr/bin/neutron-server --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/plugins/ml2/ml2_conf.ini
