# docker/openstack/heat/heat-api/root/etc/service/heat-api/run
#!/usr/bin/env sh

exec /usr/bin/heat-api
