# docker/openstack/cinder/cinder-scheduler/root/etc/service/cinder-scheduler/run
#!/usr/bin/env sh

exec /usr/bin/cinder-scheduler --config-file=/etc/cinder/cinder.conf --log-file=/dev/stdout
