# docker/openstack/neutron/neutron-dhcp-agent/root/etc/service/neutron-dhcp-agent/run
#!/usr/bin/env sh

exec /usr/bin/neutron-dhcp-agent --config-file /etc/neutron/neutron.conf --config-file /etc/neutron/dhcp_agent.ini
