# docker/openstack/neutron/neutron-linuxbridge-agent/root/etc/service/neutron-linuxbridge-agent/run
#!/usr/bin/env sh

exec /usr/bin/neutron-linuxbridge-agent --config-file /etc/neutron/plugins/ml2/ml2_conf.ini --config-dir /etc/neutron
