# docker/infrastructure/router/root/etc/service/busybox-syslogd/run
#!/bin/sh

exec /sbin/syslogd -n -O /dev/stdout
