# docker/openstack/nova/nova-spiceproxy/root/etc/service/nova-spiceproxy/run
#!/usr/bin/env sh

exec /usr/bin/nova-spicehtml5proxy
