# docker/openstack/nova/nova-api/root/etc/service/nova-api/run
#!/usr/bin/env sh

exec /usr/bin/nova-api
